# supervisor-complete.conf
# ============================================================================
# Laravel Slurry Admin API - Supervisor 完整配置
# 适用于生产环境 - 包含队列、WebSocket、监控等服务
# 配置日期: 2024-12-16
# ============================================================================

# ----------------------------------------------------------------------------
# Laravel 队列工作者配置
# ----------------------------------------------------------------------------

# 礼品卡兑换队列 - 最高优先级 (4个进程)
[program:laravel-gift-card-worker]
command=php /www/wwwroot/slurry-admin-api/artisan queue:work redis --queue=gift_card_exchange --tries=3 --timeout=300 --memory=256 --sleep=1
process_name=%(program_name)s_%(process_num)02d
numprocs=4
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/gift-card-worker.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stopwaitsecs=60
priority=100

# 预报爬虫队列 - 高优先级 (2个进程)
[program:laravel-forecast-worker]
command=php /www/wwwroot/slurry-admin-api/artisan queue:work redis --queue=forecast_crawler --tries=3 --timeout=180 --memory=256 --sleep=2
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/forecast-worker.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stopwaitsecs=60
priority=200

# 账单处理队列 - 中等优先级 (2个进程)
[program:laravel-bill-worker]
command=php /www/wwwroot/slurry-admin-api/artisan queue:work redis --queue=bill_processing --tries=3 --timeout=120 --memory=256 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/bill-worker.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stopwaitsecs=60
priority=300

# 卡密查询队列 - 中等优先级 (2个进程)
[program:laravel-card-query-worker]
command=php /www/wwwroot/slurry-admin-api/artisan queue:work redis --queue=card_query --tries=3 --timeout=300 --memory=256 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/card-query-worker.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stopwaitsecs=60
priority=400

# 微信消息队列 - 中等优先级 (2个进程)
[program:laravel-wechat-message-worker]
command=php /www/wwwroot/slurry-admin-api/artisan queue:work redis --queue=wechat-message --tries=3 --timeout=60 --memory=256 --sleep=3
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/wechat-message-worker.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stopwaitsecs=30
priority=450

# 邮件队列 - 独立处理 (1个进程)
[program:laravel-mail-worker]
command=php /www/wwwroot/slurry-admin-api/artisan queue:work redis --queue=mail --tries=3 --timeout=60 --memory=128 --sleep=5
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/mail-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stopwaitsecs=30
priority=500

# 默认队列 - 低优先级 (2个进程)
[program:laravel-default-worker]
command=php /www/wwwroot/slurry-admin-api/artisan queue:work redis --queue=high,default --tries=3 --timeout=300 --memory=256 --sleep=5
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/default-worker.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stopwaitsecs=60
priority=600

# ----------------------------------------------------------------------------
# WebSocket 服务配置
# ----------------------------------------------------------------------------

# 交易监控 WebSocket 服务器 (端口 8848)
[program:websocket-trade-monitor]
command=/bin/bash /www/wwwroot/slurry-admin-api/start-websocket-production.sh
process_name=%(program_name)s
numprocs=1
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/websocket-trade-monitor.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stopwaitsecs=10
priority=700

# ----------------------------------------------------------------------------
# 监控脚本配置
# ----------------------------------------------------------------------------

# 礼品卡日志监控服务
[program:gift-card-log-monitor]
command=/bin/bash /www/wwwroot/slurry-admin-api/monitor_gift_card_logs.sh --daemon
process_name=%(program_name)s
numprocs=1
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/gift-card-monitor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stopwaitsecs=10
priority=800

# WebSocket 监控服务
[program:websocket-monitor]
command=/bin/bash /www/wwwroot/slurry-admin-api/monitor_websocket_logs.sh
process_name=%(program_name)s
numprocs=1
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/websocket-monitor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stopwaitsecs=10
priority=900

# ----------------------------------------------------------------------------
# 计划任务调度器 (Laravel Scheduler)
# ----------------------------------------------------------------------------

# Laravel 任务调度器 - 每分钟运行
[program:laravel-scheduler]
command=/bin/bash -c "while true; do php /www/wwwroot/slurry-admin-api/artisan schedule:run --verbose --no-interaction & sleep 60; done"
process_name=%(program_name)s
numprocs=1
directory=/www/wwwroot/slurry-admin-api
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/www/wwwroot/slurry-admin-api/storage/logs/supervisor/scheduler.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stopwaitsecs=65
priority=1000

# ----------------------------------------------------------------------------
# 组配置 - 方便统一管理
# ----------------------------------------------------------------------------

# Laravel 队列工作者组
[group:laravel-workers]
programs=laravel-gift-card-worker,laravel-forecast-worker,laravel-bill-worker,laravel-card-query-worker,laravel-wechat-message-worker,laravel-mail-worker,laravel-default-worker
priority=100

# WebSocket 服务组
[group:websocket-services]
programs=websocket-trade-monitor
priority=200

# 监控服务组
[group:monitoring-services]
programs=gift-card-log-monitor,websocket-monitor
priority=300

# 系统服务组
[group:system-services]
programs=laravel-scheduler
priority=400

# ----------------------------------------------------------------------------
# 全局配置说明
# ----------------------------------------------------------------------------
# 
# 进程总数统计:
# - 礼品卡队列: 4个进程
# - 预报爬虫队列: 2个进程  
# - 账单处理队列: 2个进程
# - 卡密查询队列: 2个进程
# - 微信消息队列: 2个进程
# - 邮件队列: 1个进程
# - 默认队列: 2个进程
# - WebSocket服务: 1个进程
# - 监控服务: 2个进程
# - 调度器: 1个进程
# 总计: 19个进程
#
# 内存使用估算:
# - 队列进程: 15 × 256MB = 3.8GB
# - 其他进程: 4 × 128MB = 0.5GB
# - 总计约: 4.3GB
#
# 日志管理:
# - 所有日志统一存放在 storage/logs/supervisor/ 目录
# - 日志文件自动轮转，保留多个备份
# - 可通过 supervisorctl tail 命令查看实时日志
#
# 管理命令:
# - 启动所有服务: sudo supervisorctl start all
# - 停止所有服务: sudo supervisorctl stop all  
# - 重启所有服务: sudo supervisorctl restart all
# - 查看状态: sudo supervisorctl status
# - 重新加载配置: sudo supervisorctl reread && sudo supervisorctl update
# ============================================================================ 